# source/Build/install_packages/TLC/VS/Scripts/watchdog.timer
[Unit]
Description=ValueStor Watchdog Timer
Requires=watchdog.service

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
